# mc-bot.service
[Unit]
Description=Minecraft Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/mc-discord-bot
ExecStart=/root/botenv/bin/python /root/mc-discord-bot/bot.py
Restart=always
RestartSec=5s
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/root/mc-discord-bot/.env

[Install]
WantedBy=multi-user.target